# contrib/init/xcoind.openrc
#!/sbin/openrc-run

# backward compatibility for existing gentoo layout 
#
if [ -d "/var/lib/xcoin/.xcoin" ]; then
	XCOIND_DEFAULT_DATADIR="/var/lib/xcoin/.xcoin"
else
	XCOIND_DEFAULT_DATADIR="/var/lib/xcoind"
fi

XCOIND_CONFIGFILE=${XCOIND_CONFIGFILE:-/etc/xcoin/xcoin.conf}
XCOIND_PIDDIR=${XCOIND_PIDDIR:-/var/run/xcoind}
XCOIND_PIDFILE=${XCOIND_PIDFILE:-${XCOIND_PIDDIR}/xcoind.pid}
XCOIND_DATADIR=${XCOIND_DATADIR:-${XCOIND_DEFAULT_DATADIR}}
XCOIND_USER=${XCOIND_USER:-${XCOIN_USER:-xcoin}}
XCOIND_GROUP=${XCOIND_GROUP:-xcoin}
XCOIND_BIN=${XCOIND_BIN:-/usr/bin/xcoind}
XCOIND_NICE=${XCOIND_NICE:-${NICELEVEL:-0}}
XCOIND_OPTS="${XCOIND_OPTS:-${XCOIN_OPTS}}"

name="Xcoin Core Daemon"
description="Xcoin cryptocurrency P2P network daemon"

command="/usr/bin/xcoind"
command_args="-pid=\"${XCOIND_PIDFILE}\" \
		-conf=\"${XCOIND_CONFIGFILE}\" \
		-datadir=\"${XCOIND_DATADIR}\" \
		-daemon \
		${XCOIND_OPTS}"

required_files="${XCOIND_CONFIGFILE}"
start_stop_daemon_args="-u ${XCOIND_USER} \
			-N ${XCOIND_NICE} -w 2000"
pidfile="${XCOIND_PIDFILE}"

# The retry schedule to use when stopping the daemon. Could be either
# a timeout in seconds or multiple signal/timeout pairs (like
# "SIGKILL/180 SIGTERM/300")
retry="${XCOIND_SIGTERM_TIMEOUT}"

depend() {
	need localmount net
}

# verify
# 1) that the datadir exists and is writable (or create it)
# 2) that a directory for the pid exists and is writable
# 3) ownership and permissions on the config file
start_pre() {
	checkpath \
	-d \
	--mode 0750 \
	--owner "${XCOIND_USER}:${XCOIND_GROUP}" \
	"${XCOIND_DATADIR}"

	checkpath \
	-d \
	--mode 0755 \
	--owner "${XCOIND_USER}:${XCOIND_GROUP}" \
	"${XCOIND_PIDDIR}"

	checkpath -f \
	-o ${XCOIND_USER}:${XCOIND_GROUP} \
	-m 0660 \
	${XCOIND_CONFIGFILE}

	checkconfig || return 1
}

checkconfig()
{
	if ! grep -qs '^rpcpassword=' "${XCOIND_CONFIGFILE}" ; then
		eerror ""
		eerror "ERROR: You must set a secure rpcpassword to run xcoind."
		eerror "The setting must appear in ${XCOIND_CONFIGFILE}"
		eerror ""
		eerror "This password is security critical to securing wallets "
		eerror "and must not be the same as the rpcuser setting."
		eerror "You can generate a suitable random password using the following "
		eerror "command from the shell:"
		eerror ""
		eerror "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
		eerror ""
		eerror "It is recommended that you also set alertnotify so you are "
		eerror "notified of problems:"
		eerror ""
		eerror "ie: alertnotify=echo %%s | mail -s \"Xcoin Alert\"" \
			"admin@foo.com"
		eerror ""
		return 1
	fi
}
